# config/systemd/stratum-sha256d.service
[Unit]
Description=YiiMP Stratum SHA256D
After=network.target
Wants=network.target

[Service]
ExecStart=/usr/local/bin/stratum /etc/yiimp/stratum/sha256d.conf
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
