# config/systemd/yiimp-main.service
[Unit]
Description=YiiMP Main Loop (payments, coin updates, market data)
After=apache2.service memcached.service network.target
Wants=network.target

[Service]
User=www-data
WorkingDirectory=/var/www
ExecStart=/bin/bash /var/www/main.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
